# docker/supervisor-dcase.conf
[program:dcase]
command=/usr/bin/python3 /var/www/html/dcase/websocket.py
numprocs=1
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/dcase.log