# docker/supervisor-mongo.conf
[program:mongo]
command=/usr/bin/mongod --config /etc/mongodb.conf
numprocs=1
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/mongo.log